# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=0
loglevel=info

[program:nextjs]
command=npm start
directory=/app/templates/remote-attestation-template
user=nextjs
autostart=true
autorestart=true
stdout_logfile=/tmp/nextjs.log
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/nextjs.err
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3000",NEXT_TELEMETRY_DISABLED="1"
priority=200
startsecs=10
startretries=3

[program:python-api]
command=uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app/templates/remote-attestation-template/api
autostart=true
autorestart=true
stdout_logfile=/tmp/python.log
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/python.err
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PORT="8000",PYTHONPATH="/app/templates/remote-attestation-template/api"
user=nextjs
priority=100
startsecs=15
startretries=3

[program:bun-server]
command=bun run index.ts
directory=/app/templates/remote-attestation-template/bun-server
autostart=true
autorestart=true
stdout_logfile=/tmp/bun.log
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/bun.err
stderr_logfile_maxbytes=0
environment=PORT="8001",NODE_ENV="production",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
user=nextjs
priority=300

[group:app]
programs=nextjs,python-api,bun-server

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
