# test-supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=0
loglevel=info

[program:nextjs]
command=npm start
directory=/app/templates/remote-attestation-template
user=nextjs
autostart=true
autorestart=true
stdout_logfile=/tmp/nextjs.log
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/nextjs.err
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3000",NEXT_TELEMETRY_DISABLED="1"
priority=200
startsecs=10
startretries=3

[program:test-api]
command=python3 test-simple.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/tmp/test.log
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/test.err
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PORT="8000"
user=root
priority=100
startsecs=5
startretries=3

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
